# systemd/drive-monitor.path
[Unit]
Description = Monitoring to trigger move service
Documentation = man:systemd.path

[Path]
PathModified = /mnt/googledrive

[Install]
WantedBy = multi-user.sh